# tools/odometry/odometry.service.in
##regist service

[Unit]
Description = @SERVICE_NAME@ daemon

[Service]
ExecStart = /usr/bin/@SERVICE_NAME@ -c @CONFIG_PATH@
WorkingDirectory=@SCRIPT_DIR@

Restart = always
RestartSec = 3
Type = simple

[Install]
WantedBy = multi-user.target
