# tools/odometry/odometry.service.in
##regist service
#bash setup_service.sh /path/to/config.json

[Unit]
Description = odometry daemon

[Service]
ExecStart = /usr/bin/odometry -c @CONFIG_PATH@

#Environment="LD_LIBRARY_PATH=/path/to/cudnn:$LD_LIBRARY_PATH"

Restart = always
RestartSec = 3
Type = simple

[Install]
WantedBy = multi-user.target
