# setup/pserver.service.in
[Unit]
Description = pserver daemon

[Service]
ExecStartPre=/bin/sleep 10  # wait for EglDisplay ready
ExecStart = /bin/bash @CMAKE_INSTALL_PREFIX@/bin/pserver
ExecStartPost = /bin/sleep 3
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target