# pserver/pserver2.service.in
##regist service
#bash setup_service.sh /path/to/config.json

[Unit]
Description = pserver2 daemon

[Service]
# wait for EglDisplay ready
#ExecStartPre = /bin/sleep 3
ExecStart = /usr/bin/pserver2 -c @CONFIG_PATH@
#ExecStop = -/usr/bin/docker kill pserver

User=picam360
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/picam360/.Xauthority"

Restart = always
RestartSec = 1
Type = simple

[Install]
WantedBy = multi-user.target
