# pserver/pserver2.service.in
##regist service
#bash setup_service.sh /path/to/config.json

[Unit]
Description = pserver2 daemon

[Service]
# wait for EglDisplay ready
ExecStartPre = /bin/sleep 10
ExecStart = /usr/bin/pserver2 -c @CONFIG_PATH@
#ExecStop = -/usr/bin/docker kill pserver

Environment="DISPLAY=:0"

Restart = always
RestartSec = 3
Type = simple

[Install]
WantedBy = multi-user.target
